# services/nzbget/run
#!/bin/bash
umask 000

exec /sbin/setuser abc nzbget -s -c /config/nzbget.conf