# services/nzbget/run
#!/bin/bash
umask 000



exec /sbin/setuser abc /app/nzbget -D -c /config/nzbget.conf
PID_ID=$(cat /downloads/nzbget.lock)
while [ kill -s 0 $PID_ID]; do
sleep 500
done

rm -f /downloads/nzbget.lock