# services/nzbget/run
#!/bin/bash
umask 000

/sbin/setuser abc /app/nzbget -D -c /config/nzbget.conf
sleep 2s
exec /sbin/setuser abc tail -F /downloads/nzbget.log --pid="$(cat /downloads/nzbget.lock)"
