# services/nzbget/run
#!/bin/bash
umask 000

exec /sbin/setuser abc /app/nzbget -D -c /config/nzbget.conf