# services/tail/run
#!/bin/bash

exec /sbin/setuser abc tail -F /downloads/nzbget.log